# 5Semestre/TVS/tvs-2526-1-leic51d-coursework-2526-g04d1/cw3/systemd/tvschatd.socket
[Unit]
Description=tvschatd socket for tvschat manager

[Socket]
ListenStream=/run/isel/tvschat/request
SocketMode=0660
SocketUser=root
SocketGroup=tvsgrp

[Install]
WantedBy=sockets.target 