# ISEL/5Semestre/TVS/tvs-2526-1-leic51d-coursework-2526-g04d1/cw3/systemd/tvschatd.service
[Unit]
Description=TVS Chat Manager
Requires=tvschatd.socket

[Service]
ExecStart=/opt/isel/tvs/tvschat/manager/daemon/bin/tvschatd
Restart=always

[Install]
WantedBy=multi-user.target
Also=tvschatd.socket
