# ISEL/5Semestre/TVS/tvs-2526-1-leic51d-coursework-2526-g04d1/cw3/systemd/tvswebapp.service
[Unit]
Description=TVS Web Application
After=network.target

[Service]
Environment=NODE_PORT=54321
Type=simple
User=isel
ExecStart=/usr/bin/node /opt/isel/tvs/tvschat/webapp/tvschat.js
Restart=on-failure

